# elements/install-node-exporter/static/etc/systemd/system/node-exporter.service
[Unit]
Description = node-exporter

[Service]
User = root
ExecStart = /usr/local/bin/node_exporter --web.listen-address=:9100
Restart = always
Type = simple

[Install]
WantedBy = multi-user.target
